# begemot/extensions/crontab/crontabs/my_crontab
* * * * * php /home/users/t/timofeimih-gmail/domains/rosvesdehod.myprtf.com/yiicmd.php test livejournal.php,argo.php
