# begemot/extensions/crontab/crontabs/my_crontab
* * * * 1 php /home/users/t/timofeimih-gmail/domains/rosvesdehod.myprtf.com/modules/begemot/extensions/crontab/yiicmd.php parse argo.php
